# initscript/stap-server.in
#!/bin/bash
#
# stap-server	init.d script for the systemtap compile server
#
# chkconfig: - 00 99
# config: /etc/sysconfig/stap-server
# config: /etc/stap-server/conf.d
### BEGIN INIT INFO
# Provides: SystemTap compile-server startup
# Required-Start: $local_fs $network
# Required-Stop: $local_fs $network
# Short-Description: Start and stop SystemTap compile-server
# Description: stap-server provides remote-compilation services to systemtap users on the local network.
### END INIT INFO

# Call the installed stap-server script.
STAP_USER=stap-server @bindir@/stap-server "$@"
